# scripts/checklist-updater/checklist-updater.service
[Unit]
Description=Kamui Rossy Checklist Auto Updater
After=network.target

[Service]
Type=simple
User=hiroshi
WorkingDirectory=/mnt/c/Users/e-hac/OneDrive/デスクトップ/kamui_rossy
ExecStart=/usr/bin/python3 /mnt/c/Users/e-hac/OneDrive/デスクトップ/kamui_rossy/scripts/checklist-updater/advanced-log-analyzer.py --hours 24
Restart=always
RestartSec=3600

[Install]
WantedBy=multi-user.target